# debian/python3-jobmanager-client.prerm
#!/bin/sh

case "$1" in
    remove)
        # Stop the service
        /etc/init.d/jobmanager-client stop

        # Remove the service init script
        update-rc.d -f jobmanager-client remove
    ;;
esac
