# debian/python3-jobmanager-client.prerm
#!/usr/bin/env bash

/etc/init.d/jobmanager-client stop
update-rc.d -f jobmanager-client remove
